# luci-app-routerpilot/root/etc/uci-defaults/40_luci-routerpilot
#!/bin/sh

uci -q batch <<EOF
	del uci.routerpilot 2>/dev/null

	set routerpilot.general=general
	set routerpilot.general.port=8080
	set routerpilot.general.host=0.0.0.0
	set routerpilot.general.log_level=info
	set routerpilot.general.read_only=1
	set routerpilot.general.dry_run=0

	set routerpilot.telegram=telegram
	set routerpilot.telegram.tel_enabled=0
	set routerpilot.telegram.tel_token=""
	set routerpilot.telegram.tel_allowed_ids=""

	set routerpilot.llm=llm
	set routerpilot.llm.llm_type=simple
	set routerpilot.llm.llm_api_key=""
	set routerpilot.llm.llm_endpoint=https://api.openai.com/v1
	set routerpilot.llm.llm_model=gpt-4

	set routerpilot.schedule=schedule

	commit routerpilot
EOF

exit 0
